# upload-to-sv-machine/create-host-data.service
[Unit]
Description=create-host-data

[Service]
ExecStart=/bin/bash /home/ec2-user/create-host-data.sh

[Install]
WantedBy=multi-user.target